# .history/deploy/paike_20260313123917.service
[Unit]
Description=北清商学院排课系统
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/paike/backend
ExecStart=/usr/local/bin/gunicorn --bind 0.0.0.0:5000 --workers 2 --timeout 120 app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
